# .claude/container/crontab
# Agent container cron schedule
# All times in container timezone (set TZ in docker-compose.yml)
#
# Pattern: each job sources secrets from /tmp/agent-secrets.env (written by entrypoint.sh
# at startup via 1Password CLI injection). This avoids storing credentials in the crontab.
#
# Note: MAILTO="" suppresses cron email (no MTA in container)
MAILTO=""

# Git auto-pull (every 1 min) — keeps container in sync with GitHub pushes from laptop
# Requires /tmp/ssh-state/ssh_config to be written by entrypoint.sh at startup
* * * * * cd /home/claude/workspace && GIT_SSH_COMMAND="ssh -F /tmp/ssh-state/ssh_config" git pull --ff-only --autostash 2>&1 | head -5 >> /tmp/git-auto-pull.log

# Auth health check — 2:45 AM daily (15 min before first scheduled job)
# Validates setup-token works, alerts via Slack if it doesn't
45 2 * * * set -a; . /tmp/agent-secrets.env; set +a; bash /home/claude/workspace/.claude/scripts/auth-refresh.sh >> /home/claude/.claude/job-state/auth-refresh.log 2>&1

# Daily retro batch — 3 AM daily (2 hours before orchestrator)
# Analyzes yesterday's session issues with parallel agents + synthesis.
# Writes signal file consumed by morning orchestrator at 5 AM.
0 3 * * * set -a; . /tmp/agent-secrets.env; set +a; cd /home/claude/workspace && bash .claude/scripts/daily-retro.sh >> /home/claude/.claude/job-state/daily-retro.log 2>&1

# Morning orchestrator — 5 AM daily (VPS runs 24/7, no laptop wake needed)
# Multi-phase orchestrator: memory synthesis, briefing generation, context refresh.
0 5 * * * set -a; . /tmp/agent-secrets.env; set +a; cd /home/claude/workspace && bash .claude/scripts/morning-orchestrator.sh >> /home/claude/.claude/good-morning-state/launchd.log 2>&1

# Afternoon scan — 4 PM daily
# Mid-day context refresh and engagement scan.
0 16 * * * set -a; . /tmp/agent-secrets.env; set +a; cd /home/claude/workspace && bash .claude/scripts/afternoon-scan.sh >> /home/claude/.claude/job-state/afternoon-scan.log 2>&1

# Social post scheduler — 5x/day (example: 8,10,12,15,18)
# Uncomment and configure for your posting schedule
#0 8,10,12,15,18 * * * set -a; . /tmp/agent-secrets.env; set +a; cd /home/claude/workspace && bash Scripts/post-scheduler.sh >> /home/claude/.claude/job-state/post-scheduler.log 2>&1

# News/feed monitor — hourly during business hours
# Polls RSS feeds, drafts reaction content, pings Slack for items worth reacting to.
0 8-22 * * * set -a; . /tmp/agent-secrets.env; set +a; cd /home/claude/workspace && bash .claude/scripts/news-feed-monitor.sh >> /home/claude/.claude/job-state/news-feed-monitor.log 2>&1
